# src/ferry/data/systemd/ferry-sync.service
[Unit]
Description=ferry — sync RomM library to local ES-DE
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=__FERRY_PATH__ sync
Nice=10

# This service is timer-driven — see ferry-sync.timer. No [Install] section
# on purpose; do not enable the service directly.
